# hpsum.conf
[program:hpsum]
command=/bin/bash -c "hpsum"